# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
user=mysql
command=mysqld_safe --user=mysql --basedir=/usr/local/mysql
redirect_stderr=true

[program:tomcat]
command=tomcat run
redirect_stderr=true

